# supervisord.cloudrun.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:uvicorn]
command=uvicorn app:app --host 127.0.0.1 --port 8001 --log-level info
directory=/app/api
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=TOCG_MODEL_PATH="models/tocg.onnx",GEN_MODEL_PATH="models/gen.onnx",UPLOAD_DIR="/tmp/tryon/uploads",OUTPUT_DIR="/tmp/tryon/outputs"

[program:nextjs]
command=node server.js
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
